# packages/sysutils/sunxi-sys-utils/system.d/sunxi-wifi-load.service
[Unit]
Description=wifi driver loader
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=-/bin/sh -c "echo "123" > /proc/driver/wifi-pm/power"

[Install]
WantedBy=network.target
