# packages/network/connman/system.d/network-online.service
[Unit]
Description=Wait for network to be configured by ConnMan
Requisite=connman.service
After=connman.service
Before=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/connmand-wait-online --timeout=30
RemainAfterExit=yes

[Install]
WantedBy=network-online.target
