# packages/graphics/bcm2835-driver/system.d/fbset.timer
[Unit]
Description=Framebuffer setup

[Timer]
OnBootSec=20sec

[Install]
WantedBy=graphical.target
