# packages/libretro/retroarch/system.d/tmp-cores.mount
[Unit]
Description=Cores directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=unionfs#/storage/cores=RW:/usr/lib/libretro=RO
Where=/tmp/cores
Type=fuse
Options=cow

[Install]
WantedBy=retroarch.target
