# packages/libretro/retroarch/system.d/tmp-database.mount
[Unit]
Description=Database directory
Before=retroarch.target
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=unionfs#/storage/database=RW:/usr/share/libretro-database=RO
Where=/tmp/database
Type=fuse
Options=cow

[Install]
WantedBy=retroarch.target
