# projects/OdroidXU3/packages/xu-fan-control/system.d/fan-control.service
[Unit]
Description=Fan control

[Service]
Type=oneshot
ExecStart=/usr/bin/fan-control.sh
