# packages/devel/eglibc/system.d/nsswitch.service
[Unit]
Description=nsswitch.conf Setup
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/openelec/nsswitch-setup
RemainAfterExit=yes

[Install]
WantedBy=network.target
