# packages/x11/app/xrandr/system.d/intel-fullrange-resume.service
[Unit]
Description=Restore full range after suspend
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=-/bin/sh "/usr/lib/xorg/intel-fullrange"

[Install]
WantedBy=sleep.target
