# packages/tools/plymouth-lite/system.d/splash.service
[Unit]
Description=Show Bootsplash
DefaultDependencies=false

ConditionKernelCommandLine=!nosplash

[Service]
Type=oneshot

ExecStart=/bin/sh -c '[ -f /flash/oemsplash.png ] && /bin/ply-image /flash/oemsplash.png || /bin/ply-image /splash/splash-full.png'
ExecStartPost=/usr/bin/lsb_release
RemainAfterExit=yes
StandardOutput=tty

[Install]
WantedBy=basic.target
