# packages/lakka/RetroArch/system.d/retroarch-configure@.service
[Unit]
Description=configure Retroarch for %I

ConditionKernelCommandLine=!installer

[Service]
Type=oneshot
ExecStart=/usr/lib/libretro/retroarch-configure %I
RemainAfterExit=yes
