# projects/Gamegirl/packages/gamegirl-joypad/system.d/gamegirl-joypad.service
[Unit]
Description=Gamegirl joypad
After=gamegirl-screen.service
Requires=gamegirl-screen.service

[Service]
ExecStart=/usr/bin/gamegirl-joypad
KillMode=process
TimeoutStopSec=10
Restart=no
StartLimitInterval=0

[Install]
WantedBy=retroarch.target
