# packages/tools/rng-tools/system.d/rng-tools.service
[Unit]
Description=Hardware RNG Encryption Gatherer Daemon

[Service]
ExecStart=/usr/sbin/rngd -f

[Install]
WantedBy=basic.target
