# packages/tools/rng-tools/system.d/rng-tools.service
[Unit]
Description=RNG-Tools service

[Service]
Type=oneshot
ExecStart=/usr/sbin/rngd

[Install]
WantedBy=basic.target
