# packages/project-imx6/soc-fan/system.d/soc-fan-monitor.service
[Unit]
Description=SoC fan monitor service
ConditionPathExists=/var/run/system-id-matrix

[Service]
ExecStart=/usr/bin/soc-fan-monitor

[Install]
WantedBy=multi-user.target
