# packages/linux-drivers/gpu-aml/udev.d/81-mali.rules
ACTION=="remove", GOTO="drivers_end"

KERNEL=="amvideo", SUBSYSTEM=="module", RUN+="/sbin/modprobe mali"

LABEL="drivers_end"
