# packages/sysutils/fuse/system.d/sys-fs-fuse-connections.mount
[Unit]
Description=FUSE Control File System
DefaultDependencies=no
ConditionPathExists=/sys/fs/fuse/connections
Before=sysinit.target

[Mount]
What=fusectl
Where=/sys/fs/fuse/connections
Type=fusectl
