# packages/oem/Matrix/fan/system.d/fanmonitor.service
[Unit]
Description=Fan monitor service

[Service]
ExecStart=/usr/bin/fanmonitor

[Install]
WantedBy=multi-user.target
