# packages/linux-drivers/wetekdvb/udev.d/81-wetekdvb.rules
ACTION=="remove", GOTO="drivers_end"

KERNEL=="wetek-dvb", SUBSYSTEM=="drivers", RUN+="/sbin/modprobe wetekdvb"

LABEL="drivers_end"
